# rootfs_files/etc/runit/sv/getty-generic/run
#!/bin/sh

tty=${PWD##*-}

[ -r conf ] && . ./conf

exec chpst -P getty ${GETTY_ARGS} \
	"${tty}" "${BAUD_RATE}" "${TERM_NAME}"
