# rootfs_files/etc/runit/sv/nodm/run
#!/bin/sh
export NODM_USER="root"
export NODM_X_OPTIONS="-nolisten tcp"

nodm
